# docker/platocdp-haproxy/supervisord.conf
[supervisord]
nodaemon=true
user=haproxy
logfile=/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info

[program:postgresql]
command=haproxy -f /etc/haproxy.cfg
user=haproxy
